# deployment/web/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
user=root
stopsignal=INT
killasgroup=true
startretries=0
priority=40
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
